# etc/cron.d/zfs-auto-snapshot
PATH="/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin"

# Snapshot every 15 minutes
# Keep 6 hours of these, 24 total
*/15 * * * * root which zfs-auto-snapshot > /dev/null || exit 0 ; zfs-auto-snapshot --quiet --syslog --label=frequently --keep=24 //

# Shapshot every hour
# Keep 3 days of these, 72 total
0 * * * * root which zfs-auto-snapshot > /dev/null || exit 0 ; zfs-auto-snapshot --quiet --syslog --label=hourly --keep=72 //

# Snapshot every day
# Keep 2 weeks of these, 14 total
0 0 * * * root which zfs-auto-snapshot > /dev/null || exit 0 ; zfs-auto-snapshot --quiet --syslog --label=daily --keep=14 //

# Snapshot every week
# Keep ~2 months of these, 10 total
0 0 * * 0 root which zfs-auto-snapshot > /dev/null || exit 0 ; zfs-auto-snapshot --quiet --syslog --label=weekly --keep=10 //

# Snapshot every month
# Keep 1/2 year of these, 6 total
0 0 1 * * root which zfs-auto-snapshot > /dev/null || exit 0 ; zfs-auto-snapshot --quiet --syslog --label=monthly --keep=6 //